# crm-enrichment/crm-enrichment.service
[Unit]
Description=CRM Auto-Enrichment Service
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=/home/merce/.openclaw/workspace/crm-enrichment
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
